# aws-cloudformation/habsup.service
[Unit]
Description=The Habitat Supervisor

[Service]
ExecStart=/bin/hab sup run --permanent-peer

[Install]
WantedBy=default.target